# deploy/systemd/stellarindex-aggregator.service
# systemd unit: stellarindex-aggregator (VWAP + confidence cycle).
#
# Long-running daemon. On each tick (default 30s) walks the
# configured pair × window matrix, computes VWAP from the indexer's
# trade rows, runs the Phase 1 + Phase 2 anomaly checks, and writes
# the result + confidence score to Redis. Every hour also runs the
# baseline-refresh worker that updates `volatility_baseline_1m`.
#
# Operator wiring (after indexer is running):
#
#   sudo cp deploy/systemd/stellarindex-aggregator.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now stellarindex-aggregator.service
#
# Configuration:
#   /etc/stellarindex.toml — [aggregate] + [anomaly] sections
#
# Dependencies (must be running BEFORE the aggregator starts):
#   - postgresql.service — reads trade rows
#   - redis-server.service — writes VWAP + confidence cache
#   - stellarindex-indexer.service — supplies the trades the
#     aggregator reads. Not a strict Requires= because operators
#     occasionally stop the indexer for maintenance while leaving
#     the aggregator serving cached values; a soft `After=` is the
#     right shape.

[Unit]
Description=Stellar Index — aggregator (VWAP / TWAP / confidence per ADR-0019)
Documentation=https://github.com/Stellar-Index/StellarIndex/blob/main/docs/adr/0019-anomaly-response-and-confidence-scoring.md
After=network-online.target postgresql.service redis-server.service stellarindex-indexer.service
Wants=network-online.target
Requires=postgresql.service redis-server.service

# Restart-storm limit — [Unit] directives per systemd.unit(5). systemd
# >=230 (the Ubuntu 22.04/24.04 deploy targets) IGNORES StartLimit* if
# placed under [Service]. When more than StartLimitBurst restarts occur
# within StartLimitIntervalSec the unit enters `failed`, which the
# node_systemd_unit_state{state="failed"} alert pages on.
StartLimitIntervalSec=5min
StartLimitBurst=10

[Service]
Type=simple
User=stellarindex
Group=stellarindex
WorkingDirectory=/var/lib/stellarindex
EnvironmentFile=-/etc/default/stellarindex-ops

ExecStart=/usr/local/bin/stellarindex-aggregator -config /etc/stellarindex.toml
ExecStop=/bin/kill -TERM $MAINPID

LimitNOFILE=1048576
LimitNPROC=1048576

Restart=on-failure
RestartSec=10s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/stellarindex

[Install]
WantedBy=multi-user.target
