# deploy/systemd/stellarindex-indexer.service
# systemd unit: stellarindex-indexer (Galexie → dispatcher → decoder
# pipeline that produces canonical.Trade rows in TimescaleDB).
#
# Long-running daemon. Reads ledger meta from MinIO via the
# go-stellar-sdk's ApplyLedgerMetadata path, dispatches each event
# to its decoder per `internal/dispatcher/`, and inserts the
# resulting canonical.Trade rows. See
# `docs/architecture/ingest-pipeline.md` for the full picture.
#
# Operator wiring:
#
#   sudo cp deploy/systemd/stellarindex-indexer.service /etc/systemd/system/
#   sudo cp /etc/default/stellarindex-ops.example /etc/default/stellarindex-ops  # AWS creds
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now stellarindex-indexer.service
#
# Configuration:
#   /etc/stellarindex.toml — global config (region, sources, storage)
#   /etc/default/stellarindex-ops — AWS_ACCESS_KEY_ID + secret for MinIO read
#
# Dependencies (must be running BEFORE the indexer starts):
#   - minio.service (the source of ledger meta)
#   - postgresql.service (the trade-row destination)
#   - The `stellarindex` DB must exist + be migrated
#     (`stellarindex-migrate up`)
#
# Restart policy: on-failure with a 10s delay. A panic during a
# decoder bug should bring the binary back so a single bad ledger
# doesn't wedge the pipeline indefinitely; an actually-broken
# config (DB unreachable) loops and is observable via systemd
# state + the `stellarindex_indexer_health` Prometheus metric.

[Unit]
Description=Stellar Index — indexer (decoder pipeline → canonical trades)
Documentation=https://github.com/Stellar-Index/StellarIndex/blob/main/docs/architecture/ingest-pipeline.md
After=network-online.target minio.service postgresql.service redis-server.service
Wants=network-online.target
Requires=postgresql.service

# Restart-storm limit — [Unit] directives per systemd.unit(5). systemd
# >=230 (the Ubuntu 22.04/24.04 deploy targets) IGNORES StartLimit* if
# placed under [Service]. When more than StartLimitBurst restarts occur
# within StartLimitIntervalSec the unit enters `failed`, which the
# node_systemd_unit_state{state="failed"} alert pages on.
StartLimitIntervalSec=5min
StartLimitBurst=10

[Service]
Type=simple
User=stellarindex
Group=stellarindex
WorkingDirectory=/var/lib/stellarindex
EnvironmentFile=-/etc/default/stellarindex-ops

ExecStart=/usr/local/bin/stellarindex-indexer -config /etc/stellarindex.toml
ExecStop=/bin/kill -TERM $MAINPID

# File-descriptor limits matter — each enabled source can hold
# several network connections + DB pool + MinIO read sockets.
LimitNOFILE=1048576
LimitNPROC=1048576

# Restart on crash but with backoff to avoid hot-looping on a
# config-side break (e.g. DB unreachable). The restart-storm limit
# that marks the unit `failed` (and pages via node_systemd_unit_state)
# lives in [Unit] (StartLimit*) — see the comment there.
Restart=on-failure
RestartSec=10s

# Hardening — read-only most of the filesystem, write only to
# the working dir + the cache dir if any.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/stellarindex

[Install]
WantedBy=multi-user.target
