# deploy/systemd/stellaratlas-aggregator.service
# systemd unit: stellaratlas-aggregator (VWAP + confidence cycle).
#
# Long-running daemon. On each tick (default 30s) walks the
# configured pair × window matrix, computes VWAP from the indexer's
# trade rows, runs the Phase 1 + Phase 2 anomaly checks, and writes
# the result + confidence score to Redis. Every hour also runs the
# baseline-refresh worker that updates `volatility_baseline_1m`.
#
# Operator wiring (after indexer is running):
#
#   sudo cp deploy/systemd/stellaratlas-aggregator.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now stellaratlas-aggregator.service
#
# Configuration:
#   /etc/stellaratlas.toml — [aggregate] + [anomaly] sections
#
# Dependencies (must be running BEFORE the aggregator starts):
#   - postgresql.service — reads trade rows
#   - redis-server.service — writes VWAP + confidence cache
#   - stellaratlas-indexer.service — supplies the trades the
#     aggregator reads. Not a strict Requires= because operators
#     occasionally stop the indexer for maintenance while leaving
#     the aggregator serving cached values; a soft `After=` is the
#     right shape.

[Unit]
Description=Stellar Atlas — aggregator (VWAP / TWAP / confidence per ADR-0019)
Documentation=https://github.com/StellarAtlas/stellar-atlas/blob/main/docs/adr/0019-anomaly-response-and-confidence-scoring.md
After=network-online.target postgresql.service redis-server.service stellaratlas-indexer.service
Wants=network-online.target
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=stellaratlas
Group=stellaratlas
WorkingDirectory=/var/lib/stellaratlas
EnvironmentFile=-/etc/default/stellaratlas-ops

ExecStart=/usr/local/bin/stellaratlas-aggregator -config /etc/stellaratlas.toml
ExecStop=/bin/kill -TERM $MAINPID

LimitNOFILE=1048576
LimitNPROC=1048576

Restart=on-failure
RestartSec=10s
StartLimitInterval=5min
StartLimitBurst=10

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/stellaratlas

[Install]
WantedBy=multi-user.target
