# deploy/systemd/stellarindex-completeness.service
[Unit]
Description=Stellar Index — incremental completeness verify (ADR-0033: substrate + recognition + projection vs the ClickHouse lake, since last watermark)
After=clickhouse-server.service postgresql.service network-online.target
Wants=clickhouse-server.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/completeness-incremental.sh
# Read-only verify, but the projection re-derive can be CH-heavy — keep it gentle
# so it never disturbs serving. Incremental scope keeps each run to minutes.
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Runs as root: the script sources /etc/default/stellarindex-ops
# (Postgres DSN + credentials) directly, which is root-readable only.
# It is read-only on served data and writes nothing to local disk —
# verdicts go to Postgres over the network — so no ReadWritePaths are
# required. Hardening mirrors the archive-completeness root oneshot.
User=root
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
