# configs/ansible/roles/archival-node/templates/systemd/ch-supply.service.j2
# systemd unit: supply_flows gap-fill (SEP-41 supply).
#
# Rendered by Ansible. Runs run-ch-supply.sh: defensive forward-gap-fill of
# stellar.supply_flows from the certified ClickHouse lake — the per-token
# mint/burn/clawback set that backs /v1/assets SEP-41 supply (TokenSupply() sums
# it FINAL on-demand; the 10k+ Soroban tokens can't sit on an LCM observer
# watch-list). supply_flows is normally kept current by the indexer's live
# decode-at-ingest path, so this is usually a no-op; it catches the live writer
# up if it ever falls behind.

[Unit]
Description=Stellar Index — supply_flows gap-fill (SEP-41 supply, daily)
Documentation=https://github.com/StellarIndex/stellar-index/blob/main/docs/architecture/clickhouse-supply-from-ch.md
After=network-online.target postgresql.service clickhouse-server.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/ch-supply
EnvironmentFile=-/etc/default/stellarindex

Environment=CONFIG_PATH=/etc/stellarindex.toml
Environment=CH_ADDR=127.0.0.1:9300
Environment=CHSUPPLY_CHUNK=25000

ExecStart=/usr/local/sbin/run-ch-supply.sh

# De-privileged (CS-118, audit-2026-06-30) — the stellarindex user is
# created by 04-users.yml (CS-119). Talks Postgres + ClickHouse over the
# loopback only; writes nothing to disk.
User=stellarindex
Group=stellarindex
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# First run catches up [62.94M, tip] (the seed gap) + a full aggregate; cap
# generously. Steady-state runs seed one day + re-aggregate in minutes.
TimeoutStartSec=240min

[Install]
WantedBy=multi-user.target
