# configs/ansible/roles/archival-node/templates/systemd/ch-supply.timer.j2
# systemd timer: daily supply_flows gap-fill (SEP-41 supply).
#
# Rendered by Ansible. Companion to ch-supply.service. Defensive backstop that
# keeps stellar.supply_flows complete to tip so /v1/assets SEP-41 supply (summed
# from it on-demand) tracks the network even if the live decode-at-ingest writer
# ever falls behind. Normally a no-op.
#
# Time slot: 06:00 UTC + jitter, after compute-completeness (05:30). The seed is
# chunked + memory-guarded so it never collides with the morning oneshots'
# ClickHouse load. Persistent=true: the wrapper's chunked seed absorbs a
# post-outage backlog safely.

[Unit]
Description=Daily supply_flows gap-fill (SEP-41 lake-derived supply)
Documentation=https://github.com/Stellar-Index/StellarIndex/blob/main/docs/architecture/clickhouse-supply-from-ch.md

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
RandomizedDelaySec=300
Persistent=true
Unit=ch-supply.service

[Install]
WantedBy=timers.target
