# configs/ansible/roles/archival-node/templates/systemd/stellaratlas-indexer.service.j2
[Unit]
Description=Stellar Atlas — indexer (Galexie → dispatcher → trades hypertable)
After=network-online.target postgresql@{{ postgres_version }}-main.service minio.service galexie.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/stellaratlas
EnvironmentFile=/etc/default/stellaratlas
ExecStart=/usr/local/bin/stellaratlas-indexer -config /etc/stellaratlas.toml
ExecStop=/bin/kill -TERM $MAINPID

LimitNOFILE=1048576
Restart=on-failure
RestartSec=10s

NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
