# configs/ansible/roles/archival-node/templates/systemd/stellarindex-api.service.j2
# Rendered by Ansible (archival-node role).
#
# CS-118 (audit-2026-06-30): the INTERNET-FACING API no longer runs
# as root — it runs de-privileged as the stellarindex service user
# with the hardening block from the in-repo unit at
# deploy/systemd/stellarindex-api.service (which was drifting unused
# while this template ran the daemon as root). The stellarindex user
# is created by 04-users.yml (CS-119). EnvironmentFile is read by
# systemd (PID 1), not the service user; it stays 0640
# root:stellarindex. Caddy terminates TLS and proxies to :3000, so
# no privileged port bind is needed.

[Unit]
Description=Stellar Index — REST + SSE API server
After=network-online.target postgresql@{{ postgres_version }}-main.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=stellarindex
Group=stellarindex
WorkingDirectory=/var/lib/stellarindex
EnvironmentFile=/etc/default/stellarindex
ExecStart=/usr/local/bin/stellarindex-api -config /etc/stellarindex.toml
ExecStop=/bin/kill -TERM $MAINPID

LimitNOFILE=1048576
Restart=on-failure
RestartSec=10s

# Hardening — read-only most of the filesystem, write only the
# working dir. Mirrors deploy/systemd/stellarindex-api.service.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/stellarindex

[Install]
WantedBy=multi-user.target
