# configs/ansible/roles/archival-node/templates/systemd/stellarindex-aggregator.service.j2
[Unit]
Description=Stellar Index — aggregator (VWAP/TWAP/divergence/freeze/supply tick)
After=network-online.target postgresql@{{ postgres_version }}-main.service redis-server.service stellarindex-indexer.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/stellarindex
EnvironmentFile=/etc/default/stellarindex
ExecStart=/usr/local/bin/stellarindex-aggregator -config /etc/stellarindex.toml
ExecStop=/bin/kill -TERM $MAINPID

LimitNOFILE=1048576
Restart=on-failure
RestartSec=10s

NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
