# configs/healthchecks/ratesengine-sla-probe.service
[Unit]
Description=Rates Engine SLA probe (RFP latency + freshness pass/fail)
Documentation=https://github.com/RatesEngine/rates-engine/blob/main/configs/healthchecks/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/ratesengine/healthchecks/sla-probe.sh
EnvironmentFile=-/etc/default/ratesengine-healthchecks
DynamicUser=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
