# configs/ansible/roles/loki/templates/loki.service.j2
{# Render-time — managed by ansible (loki role). #}

[Unit]
Description=Grafana Loki — log aggregation
Documentation=https://grafana.com/docs/loki/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=loki
Group=loki

# F-1286 (codex audit-2026-05-13): MinIO credentials must come
# from `/etc/default/loki` DIRECTLY as `AWS_ACCESS_KEY_ID=` and
# `AWS_SECRET_ACCESS_KEY=`. Earlier templates used
# `Environment=AWS_ACCESS_KEY_ID=${STELLARATLAS_S3_ACCESS_KEY}`
# attempting to reuse the galexie convention's env-var names,
# but systemd `Environment=` does NOT perform `$` expansion —
# the literal `${STELLARATLAS_S3_ACCESS_KEY}` string would be
# passed to Loki as the access key, failing S3 chunk writes
# while the process otherwise looked healthy.
#
# Operators put both into /etc/default/loki:
#
#   AWS_ACCESS_KEY_ID=<value>
#   AWS_SECRET_ACCESS_KEY=<value>
#
# (Galexie's separate /etc/default/galexie can still use the
# STELLARATLAS_S3_* convention — these are two different
# EnvironmentFile sources.)
EnvironmentFile=-/etc/default/loki

ExecStart=/usr/local/bin/loki \
  -config.file=/etc/loki/loki-config.yaml

KillMode=process
TimeoutStopSec=30
LimitNOFILE=65536

# Hardening
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
